# pocketbase-server/deploy/isivolt-pocketbase-backup.timer
[Unit]
Description=Copia diaria de ISIVOLT PocketBase

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=10m
Unit=isivolt-pocketbase-backup.service

[Install]
WantedBy=timers.target
